# deploy/credential-gateway.service
[Unit]
Description=credential-gateway (local credential injection proxy)
Documentation=https://github.com/SHUKE-LABS/credential-gateway
# All listeners bind loopback only (see README), so there is no network
# interface to wait for at boot — no network-online.target / tailscaled ordering.

[Service]
Type=simple
# Runs as root: the config is 0600 root:root and is read directly via
# os.ReadFile (not an EnvironmentFile), so a DynamicUser's random UID could
# not read it. -config is passed explicitly rather than relying on the
# $HOME-based search path, which is not predictable under systemd with no User=.
# ExecStartPre validates the config before every start/restart — deploy and the
# operator's manual `systemctl restart` alike — so a bad config fails with a
# clean, greppable journal line ("failed to load config") instead of a bind or
# parse error mid-startup, and no port is ever bound. -validate only reads the
# config and exits, so it runs fine under the ProtectSystem=strict sandbox below.
ExecStartPre=/usr/local/bin/credential-gateway -validate -config /etc/credential-gateway/config.yaml
ExecStart=/usr/local/bin/credential-gateway -config /etc/credential-gateway/config.yaml
Restart=always
RestartSec=2
# Listeners drain within 10s on SIGTERM (main.go); give a hair more before SIGKILL.
TimeoutStopSec=15

# --- Hardening -----------------------------------------------------------
# Root here only needs to read one config file and open TCP sockets; nothing
# else. Stateless service, so no StateDirectory / no writable paths.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
# Only IP sockets are needed (loopback listeners + outbound TCP to upstreams).
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
