# deploy/credential-gateway-admin.service
[Unit]
Description=credential-gateway-admin (loopback config editor UI)
Documentation=https://github.com/SHUKE-LABS/credential-gateway
# Binds 127.0.0.1 only (reached via SSH tunnel), so there is no network
# interface to wait for at boot.

[Service]
Type=simple
# Runs as a dedicated non-root user with a stable UID: config.yaml stays
# 0600 root:root and cg-admin reaches it via a POSIX ACL granted at install
# (see remote-install.sh). NOT DynamicUser — the ACL grant needs a stable UID.
User=cg-admin
Group=cg-admin
ExecStart=/usr/local/bin/credential-gateway-admin -config /etc/credential-gateway/config.yaml -port 8099
Restart=always
RestartSec=2
TimeoutStopSec=15

# --- Hardening -----------------------------------------------------------
NoNewPrivileges=yes
ProtectSystem=strict
# ProtectSystem=strict makes /etc read-only; grant write to exactly the one
# file the UI edits, nothing else under /etc.
ReadWritePaths=/etc/credential-gateway/config.yaml
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
# Only IP sockets are needed (loopback HTTP listener).
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
